# systemd/vpn-control.service
[Unit]
Description=VPN Control Plane - WireGuard Admin Dashboard
After=network.target wireguard.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/opt/vpn-control
ExecStart=/opt/vpn-control/venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 8000
Restart=always
RestartSec=5

# Security hardening (while still allowing WireGuard management)
NoNewPrivileges=false
ProtectSystem=false
ProtectHome=read-only
PrivateTmp=true

# Environment
Environment="VPN_ENDPOINT=vpn.example.com:51820"
Environment="SESSION_SECRET=change-this-use-openssl-rand-hex-32"

[Install]
WantedBy=multi-user.target
